# projects/mistral/backend/cron/arki-mergeconf.cron
SHELL=/bin/bash

0 * * * * set -a && source /etc/rapydo-environment && arki-mergeconf --extra /arkimet/datasets/* > /arkimet/config/arkimet.conf.tmp && mv /arkimet/config/arkimet.conf.tmp /arkimet/config/arkimet.conf

