# projects/mistral/backend/cron/arki-mergeconf.cron
SHELL=/bin/bash

0 * * * * echo "merging... " && /bin/bash /code/mistral/cron/arki-mergeconf.sh >> /logs/arki-mergeconf.log 2>&1
