# projects/mistral/backend/cron/arki-mergeconf.cron
0 * * * * arki-mergeconf --extra /arkimet/datasets/* > /arkimet/config/arkimet.conf.tmp && mv /arkimet/config/arkimet.conf.tmp /arkimet/config/arkimet.conf

