# projects/mistral/backend/cron/arki-mergeconf.cron
SHELL=/bin/bash

0 * * * * /bin/bash /code/mistral/cron/arki-mergeconf.cron
