# projects/mistral/backend/cron/create_dballe_json_summary.cron
*/10 * * * * PYTHONPATH=/code python3 /scripts/create_dballe_summary.py >> /var/log/cron.log 2>&1

