# projects/mistral/backend/cron/create_dballe_json_summary.cron
SHELL=/bin/bash

*/10 * * * * source /etc/rapydo-environment && PYTHONPATH=/code python3 /scripts/create_dballe_summary.py >> /var/log/cron.log 2>&1
