# projects/mistral/backend/cron/start_dballe_aggr_cleanup.cron
SHELL=/bin/bash

5 2 * * * set -a && source /etc/rapydo-environment && python3 /scripts/dballe_aggr_cleanup.py >> /var/log/cron.log 2>&1
