# projects/mistral/backend/cron/start_dballe2arkimet.cron
SHELL=/bin/bash

15 2 * * * set -a && source /etc/rapydo-environment && /scripts/start_dballe2arkimet.sh >> /var/log/cron.log 2>&1
