# projects/mistral/backend/cron/start_dballe2arkimet.cron
15 2 * * * /scripts/start_dballe2arkimet.sh >> /var/log/cron.log 2>&1

